# udp-traffic-gen.service
[Unit]
Description=UDP Traffic Generator
After=network.target

[Service]
Type=simple
User=labuser
WorkingDirectory=/home/labuser/udp-traffic-gen
ExecStart=/home/labuser/udp-traffic-gen/venv/bin/python3 /home/labuser/udp-traffic-gen/app.py
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
